# config/.config/systemd/user/offlineimap-oneshot.service.d/override.conf
[Service]
ExecStop=/usr/bin/emacsclient --eval "(update-mail-in-server)"