# config/.config/systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m